# etc/udev/rules.d/99-usb_mouse_add.rules
KERNEL=="hiddev0",SUBSYSTEM=="usbmisc",ACTION=="add",RUN+="/home/albamc/env/disable_touchpad.sh"
