# etc/udev/rules.d/99-usb_mouse_remove.rules
KERNEL=="hiddev0",SUBSYSTEM=="usbmisc",ACTION=="remove",RUN+="/home/albamc/env/enable_touchpad.sh"
